# scripts/jumpin.conf
description     "launch jumpin"

start on startup

stop on runlevel [!023456]

task
exec /etc/init.d/jumpin start